# deploy/templates/loop-layer2.timer
[Unit]
Description=Fire ops-loop Layer 2 (Orient) floor daily at 12:00 Europe/Paris
Documentation=https://github.com/Bubble-invest/bubble-ops-loop

[Timer]
# +2h backup offset ({{OPERATOR}} msg 3911, 2026-06-06): the live /loop owns each
# layer at its min-time (L2); this safety-net cron fires 2h later to catch up
# on a STALE loop without racing or replacing it.
Unit=loop-layer@2.service
# 12:00 Europe/Paris — systemd handles CET/CEST automatically.
OnCalendar=*-*-* 14:00:00 Europe/Paris
# If the box was off/asleep at fire time, run as soon as it's back.
Persistent=true
# Small deterministic jitter so it doesn't collide with other top-of-hour jobs.
RandomizedDelaySec=120

[Install]
WantedBy=timers.target
