# deploy/templates/loop-layer1.timer
[Unit]
Description=Fire ops-loop Layer 1 (Observe) floor daily at 07:00 Europe/Paris
Documentation=https://github.com/Bubble-invest/bubble-ops-loop

[Timer]
# +2h backup offset (Joris msg 3911, 2026-06-06): the live /loop owns each
# layer at its min-time (L1); this safety-net cron fires 2h later to catch up
# on a STALE loop without racing or replacing it.
Unit=loop-layer@1.service
# 07:00 Europe/Paris — systemd handles CET/CEST automatically.
OnCalendar=*-*-* 09:00:00 Europe/Paris
# If the box was off/asleep at fire time, run as soon as it's back.
Persistent=true
# Small deterministic jitter so it doesn't collide with other top-of-hour jobs.
RandomizedDelaySec=120

[Install]
WantedBy=timers.target
