# deploy/templates/secrets-tmp-sweep.service
[Unit]
Description=Bubble Security — plaintext-secret /tmp leak sweep (source-leak detector)
Documentation=https://github.com/Bubble-Finance-Inc/bubble-ops-loop
After=network-online.target
Wants=network-online.target
ConditionPathExists=/run/claude-agent/env

[Service]
Type=oneshot
User=claude
Group=claude
ExecStart=/home/claude/scripts/secrets-tmp-sweep.sh
# Grant read-only DAC bypass so the sweep can SEE root-owned 0600 secret leftovers
# (it only stats names+perms, never reads contents). NoNewPrivileges must be false
# for ambient caps to apply.
NoNewPrivileges=false
AmbientCapabilities=CAP_DAC_READ_SEARCH
CapabilityBoundingSet=CAP_DAC_READ_SEARCH
# IMPORTANT: NO PrivateTmp — this sweep MUST see the real /tmp to detect leaks.
ProtectSystem=strict
ReadWritePaths=/var/log/bubble-security
ReadOnlyPaths=/tmp /var/tmp /home/claude /run/claude-agent/env
TimeoutStartSec=120
StandardOutput=journal
StandardError=journal
SyslogIdentifier=secrets-tmp-sweep

[Install]
WantedBy=multi-user.target
