# deploy/templates/secrets-tmp-sweep.timer
[Unit]
Description=Run plaintext-secret /tmp leak sweep every 30 minutes
Documentation=https://github.com/Bubble-Finance-Inc/bubble-ops-loop

[Timer]
OnBootSec=5min
OnUnitActiveSec=30min
Persistent=true
RandomizedDelaySec=60
Unit=secrets-tmp-sweep.service

[Install]
WantedBy=timers.target
