# deploy/templates/onfailure-dropins/bubble-deploy-full.service.d/override.conf
# WS-B alert wiring: generic OnFailure -> Telegram alert (cron-failure-alert.sh).
# See deploy/templates/onfailure-dropins/README.md.
[Unit]
OnFailure=cron-failure-alert@%n.service
